# vboxguest/x86_64/usr/lib/systemd/system/vboxguest.service
[Unit]
Description=Initialize VirtualBox Guest Additions
After=local-fs.target
Before=network.service
Wants=local-fs.target network.service

[Service]
Type=oneshot
ExecStart=/etc/cernvm/service.d/vboxguest
RemainAfterExit=yes
TimeoutSec=0

[Install]
WantedBy=multi-user.target

