# copilot-jobmanager-generic/files
/etc/copilot/*
/usr/bin/*
/etc/init.d/copilot-jobmanager-generic
