# vboxguest/files
/usr/bin/*
/usr/sbin/*
/usr/share/*
/usr/lib64/*
/sbin/*
/etc/init.d/*
