# vboxguest/files
/usr/bin/*
/usr/sbin/*
/usr/share/*
/usr/lib64/*
/sbin/*
/etc/cernvm/service.d/*
/usr/lib/systemd/system/*
