# cernvm-config/files
/.ucernvm_bootstrap
/.ucernvm_pinfiles
/etc/fstab
/etc/workqueue.conf
/etc/*.cernvm
/etc/cloud/*
/etc/condor/*
/etc/cvmfs/*
/etc/init.d/*
/etc/pam.d/*
/etc/squid/*
/etc/polkit-1/*
/etc/profile.d/*
/etc/ssh/*
/etc/X11/*
/etc/sysconfig/*
/etc/tune-profiles/*
/etc/cernvm
/etc/dhcp/dhclient.d/*
/etc/modprobe.d/*
/etc/selinux/*
/etc/xrdp/*
/usr/share/Modules/modulefiles/*
/etc/grid-security/*
/usr/bin/*
