# cernvm-online-guest/files
/etc/cernvm/keys/*
/usr/sbin/cernvm-online.sh
/usr/lib/systemd/system/cernvm-online.service
