# vboxguest/x86_64/usr/lib/systemd/system/ntpd.service.d/10-cernvm.conf
[Unit]
ConditionPathExists=!/etc/cernvm/tools.d/virtualbox
After=vboxguest.service
