# copilot-agent/files
/usr/bin/copilot-agent
/etc/copilot/*
/etc/init.d/copilot-agent
