# cernvm-release/files
/etc/init.d/cernvm-release
/etc/yum.repos.d/*
