# cernvm-release/files
/etc/cernvm/service.d/cernvm-release
/usr/lib/systemd/system/cernvm-release.service
