# cernvm-update/files
/usr/bin/*
/var/lib/cernvm-update
/etc/cernvm/cernvm.d/*
/etc/cron.daily/*
/usr/share/cernvm-update
